# network-monitor-syslog.service
[Unit]
Description=Network Monitor Syslog UDP Listener Service
After=network.target network-monitor-web.service # Ensure web service is potentially up for DB access
Requires=network-monitor-web.service # Or adjust depending on dependencies

[Service]
User=root # Required for privileged port 514, or change port and user
Group=nogroup
WorkingDirectory=/opt/openwrt-loghub # Adjust path as needed
Environment="PATH=/opt/openwrt-loghub/venv/bin"
Environment="FLASK_APP=wsgi.py"
Environment="FLASK_CONFIG=production"

# Note: Running Flask development server is NOT recommended for production.
# Consider a more robust UDP server implementation if needed.
# This uses the custom Flask CLI command.
# Adjust port if not using 514.
ExecStart=/opt/openwrt-loghub/venv/bin/flask run-syslog --port 514

Restart=always

[Install]
WantedBy=multi-user.target 