# Network_Monitor/network-monitor-web.service
[Unit]
Description=Network Monitor Web Backend Service (Gunicorn)
After=network.target

[Service]
User=netmonitor # Corrected
Group=netmonitor # Corrected
WorkingDirectory=/opt/network-monitor # Corrected
Environment="PATH=/opt/network-monitor/backend/venv/bin" # Corrected
Environment="FLASK_APP=wsgi.py" # Keep wsgi for gunicorn
Environment="FLASK_CONFIG=production"
EnvironmentFile=/opt/network-monitor/.env # Explicitly load .env

ExecStart=/opt/network-monitor/backend/venv/bin/gunicorn \
    --workers 3 \
    --bind unix:/tmp/network-monitor.sock \
    # -m 0660 # Leave commented out as default permissions work reliably
    wsgi:app

Restart=always

[Install]
WantedBy=multi-user.target 