# Network_Monitor/network-monitor-syslog.service
[Unit]
Description=Network Monitor Syslog UDP Listener Service
After=network.target network-monitor-web.service # Ensure web service is potentially up for DB access
Requires=network-monitor-web.service # Or adjust depending on dependencies

[Service]
User=netmonitor
Group=netmonitor
WorkingDirectory=/opt/network-monitor # Corrected
Environment="PATH=/opt/network-monitor/backend/venv/bin" # Corrected
Environment="FLASK_APP=backend.app:create_app" # Corrected
Environment="FLASK_CONFIG=production"
EnvironmentFile=/opt/network-monitor/.env # Ensure .env is loaded

# Note: Running Flask development server is NOT recommended for production.
# Consider a more robust UDP server implementation if needed.
# This uses the custom Flask CLI command 'run-syslog'
# Port is determined from SYSLOG_UDP_PORT in .env (loaded via EnvironmentFile)
ExecStart=/opt/network-monitor/backend/venv/bin/flask run-syslog # Corrected path, removed hardcoded port

Restart=always

[Install]
WantedBy=multi-user.target 