# Network_Monitor/network-monitor-web.service
[Unit]
Description=Network Monitor Web Backend Service (Gunicorn)
After=network.target

[Service]
User=pi # Or the user you are running the app as
Group=www-data # Or the group for web access
WorkingDirectory=/opt/openwrt-loghub # Adjust path as needed
Environment="PATH=/opt/openwrt-loghub/venv/bin"
Environment="FLASK_APP=wsgi.py"
Environment="FLASK_CONFIG=production"
# Ensure your .env file is in the WorkingDirectory or specify path:
# Environment="DOTENV_PATH=/opt/openwrt-loghub/.env"

ExecStart=/opt/openwrt-loghub/venv/bin/gunicorn \
    --workers 3 \
    --bind unix:/tmp/network-monitor.sock \
    -m 007 \
    wsgi:app

Restart=always

[Install]
WantedBy=multi-user.target 